# backend/nginx/crontab
# Reload nginx every six hours
*/05 * * * * nginx -s reload
